# config/php/fpm.conf
[program:fpm]
command=/usr/local/sbin/php-fpm
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
process_name=%(program_name)s_%(process_num)02d
